# 06-storage-optimization/etc/systemd/system/fstrim-ssd-storage.service
[Unit]
Description=Custom periodic TRIM of all SSDs

[Service]
Type=oneshot
ExecStart=/usr/sbin/fstrim -av

[Install]
WantedBy=timers.target

